# scripts/com.testwrightai.watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  TestwrightAI Watchdog Sistem Servisi
  Kurulum:
    sudo cp scripts/com.testwrightai.watchdog.plist /Library/LaunchDaemons/
    sudo launchctl load /Library/LaunchDaemons/com.testwrightai.watchdog.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.testwrightai.watchdog</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/yasin_bulgan/Desktop/BGTS_Test_Donusum/scripts/twai_watchdog.sh</string>
    </array>

    <key>KeepAlive</key>
    <true/>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/twai-watchdog.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/twai-watchdog-error.log</string>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
